# pulseha.service
[Unit]
After=network.target
Requires=network.target
Description=PulseHA Daemon

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/sbin/pulse
ExecReload=/bin/kill -SIGUSR2 $MAINPID
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
